# .config/systemd/user/onedrive-sync.timer
[Unit]
Description=Syncs onedrive every 15 min

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
